# nss/overlay/target/linux/qualcommax/ipq50xx/base-files/etc/uci-defaults/00-nss-manual
#!/bin/sh
# NSS hardware offload engages at boot via the kmodloader autoload
# (/etc/modules.d/33-qca-nss-ecm loads ecm, which pulls in qca-nss-drv and boots
# the NSS core via the reset-pulse patch); redirect/accel are set from
# /etc/rc.local. Disable the package init scripts so they don't ALSO load the
# modules (double-load / wrong ordering).
#
# To back the offload out entirely (slowpath software routing always works):
#   rm /etc/modules.d/32-qca-nss-drv /etc/modules.d/33-qca-nss-ecm \
#      /etc/modules.d/51-qca-nss-drv-pppoe        # keep 31-qca-nss-dp
#   reboot
# Removing only 33-qca-nss-ecm is NOT enough: qca-nss-drv alone switches
# nss-dp onto the NSS data plane, and under a fw/driver version mismatch that
# kills the whole LAN with no ecm involved (see docs/nss-offload.md
# "Troubleshooting"). Never rmmod qca-nss-drv on a live system - it resets
# the Q6 remoteproc and crashes the box; remove the autoloads and reboot.
/etc/init.d/qca-nss-ecm disable 2>/dev/null
/etc/init.d/qca-nss-drv disable 2>/dev/null
exit 0
