# nss/overlay/target/linux/qualcommax/ipq50xx/base-files/etc/hotplug.d/iface/99-heal-wan-conduit
# Workaround for issue #5: on a clean boot the standalone wan port (conduit
# eth0) passes no CPU traffic. DSA programs the tag_8021q VLANs while the
# port's upstream is still the default conduit (eth1/P5); netifd's later
# `wan -> eth0` conduit move only rewrites the port matrix, so the VLAN state
# stays wrong for the P4/eth0 CPU path. Cycling the port through a bridge
# join/leave re-programs the tag_8021q state with the conduit already correct
# and fully heals the port until the next reboot.
#
# Runs once per boot on the first ifup of wan. Remove once
# an8855_port_change_conduit() re-programs the tag_8021q state itself.
[ "$ACTION" = ifup ] && [ "$INTERFACE" = wan ] || return 0
[ -e /tmp/.wan-conduit-healed ] && return 0
ip link show br-lan >/dev/null 2>&1 || return 0
touch /tmp/.wan-conduit-healed
(
	sleep 2
	ip link set wan down
	ip link set wan master br-lan 2>/dev/null
	ip link set wan up
	sleep 1
	ip link set wan nomaster
	ifup wan
) &
